# docker/supervisord.conf
[supervisord]
nodaemon=true
logfile=/data/logs/supervisord.log
logfile_maxbytes=10MB
loglevel=info

[program:postgres]
command=/usr/lib/postgresql/16/bin/postgres -D /var/lib/postgresql/data
user=postgres
autostart=true
autorestart=true
priority=10
stdout_logfile=/data/logs/postgres.log
stderr_logfile=/data/logs/postgres.log

[program:redis]
command=redis-server /etc/redis/redis-standalone.conf
autostart=true
autorestart=true
priority=20
stdout_logfile=/data/logs/redis.log
stderr_logfile=/data/logs/redis.log

[program:node]
command=/data/start-node.sh
autostart=false
autorestart=true
priority=30
startsecs=5
stdout_logfile=/data/logs/node.log
stderr_logfile=/data/logs/node.log
